# deploy/systemd/cleanup.service
[Unit]
Description=Qdrant Vector Cleanup Service
Documentation=https://github.com/yourusername/document-embedding-project
After=network.target

[Service]
Type=oneshot
User=embeduser
Group=embeduser
WorkingDirectory=/opt/vecpipe

# Environment
Environment="PYTHONPATH=/opt/vecpipe"
Environment="QDRANT_HOST=192.168.1.173"
Environment="QDRANT_PORT=6333"

# Execute cleanup script
ExecStart=/usr/bin/python3 /opt/vecpipe/cleanup.py --file-list /var/embeddings/filelist.null

# Resource limits
MemoryLimit=2G
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vecpipe-cleanup

# Restart policy - don't restart oneshot services
Restart=no

[Install]
WantedBy=multi-user.target