# deploy/systemd/cleanup.timer
[Unit]
Description=Daily Qdrant Vector Cleanup Timer
Documentation=https://github.com/yourusername/document-embedding-project
Requires=cleanup.service

[Timer]
# Run daily at 2 AM
OnCalendar=daily
AccuracySec=1h

# Run on startup if missed
Persistent=true

# Random delay up to 30 minutes to avoid load spikes
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target