# deploy/systemd/extract.service
[Unit]
Description=Document Extraction Service
After=manifest.service

[Service]
Type=oneshot
WorkingDirectory=/root/document-embedding-project
ExecStart=/usr/bin/python3 vecpipe/extract_chunks.py \
    --input /var/embeddings/filelist.null \
    --output /opt/vecpipe/extract \
    --batch-size 100
User=root
StandardOutput=journal
StandardError=journal
TimeoutStartSec=3600
EnvironmentFile=/etc/vecpipe/config.env

[Install]
WantedBy=multi-user.target