# deploy/systemd/ingest.service
[Unit]
Description=Qdrant Ingestion Service
After=embed.service

[Service]
Type=oneshot
WorkingDirectory=/root/document-embedding-project
ExecStart=/usr/bin/python3 vecpipe/ingest_qdrant.py \
    --input /var/embeddings/ingest \
    --loaded /var/embeddings/loaded \
    --rejects /var/embeddings/rejects
User=root
StandardOutput=journal
StandardError=journal
TimeoutStartSec=1800
EnvironmentFile=/etc/vecpipe/config.env

[Install]
WantedBy=multi-user.target