# deploy/systemd/search-api.service
[Unit]
Description=Document Search API Service
After=network.target

[Service]
Type=exec
WorkingDirectory=/root/document-embedding-project
ExecStart=/usr/bin/python3 vecpipe/search_api.py
User=root
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
EnvironmentFile=/etc/vecpipe/config.env

[Install]
WantedBy=multi-user.target